# recipes-afero/af-conn/files/af-conn.service
[Unit]
Description = Afero Connectivity
After = attrd.service

[Service]
Type=simple
ExecStart=/usr/lib/af-conn/af_conn_watcher
Restart=always
RestartSec=5

[Install]
WantedBy = multi-user.target
